# systemd/th3cl4w-map.service
[Unit]
Description=th3cl4w Map Server (port 8083)
PartOf=th3cl4w.target
After=th3cl4w-camera.service th3cl4w-location.service

[Service]
Type=simple
WorkingDirectory=/home/sarah/.openclaw/workspace/th3cl4w
EnvironmentFile=/home/sarah/.openclaw/workspace/th3cl4w/.env
ExecStart=/home/linuxbrew/.linuxbrew/bin/python3.12 web/map_server.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=th3cl4w.target
